# deploy/ehds-cpu-watchdog.timer
[Unit]
Description=Run the EHDS CPU runaway watchdog every 10 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
